# src/trashexpiry.timer
[Unit]
Description=Discard old files from Trash daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
