# src/trashexpiry.service
[Unit]
Description=Discard old files from Trash

[Service]
Type=oneshot
ExecStart=trashexpiry
